# etc/systemd/system/pwn.college.imagepuller.timer
[Unit]
Description=Timer to run pwn.college imagepuller service

[Timer]
OnCalendar=*-*-* *:*:00

[Install]
WantedBy=timers.target
