# etc/systemd/system/pwn.college.cloud.backup.timer
[Unit]
Description=Timer to run pwn.college cloud backup service

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target
