# etc/systemd/system/pwn.college.backup.timer
[Unit]
Description=Timer to run pwn.college local backup service

[Timer]
OnCalendar=hourly

[Install]
WantedBy=timers.target
