# etc/systemd/system/pwn.college.cachewarmer.timer
[Unit]
Description=Timer to run pwn.college cachewarmer service

[Timer]
OnCalendar=*-*-* *:*:00

[Install]
WantedBy=timers.target
